# server/src/pkg/moustack-server.service
[Unit]
Description=${project.name}
Documentation=${project.url}

[Service]
ExecStart=/usr/bin/moustack-server

[Install]
WantedBy=multi-user.target
